# systemd/qubes-u2fproxy.service
[Unit]
Description=Default instance of qubes-u2f-proxy

[Service]
Type=oneshot
ExecStart=/bin/true

[Install]
Also=qubes-u2fproxy@sys-usb.service
