# systemd/qubes-ctapproxy.service
[Unit]
Description=Default instance of qubes-ctap-proxy

[Service]
Type=oneshot
ExecStart=/bin/true

[Install]
Also=qubes-ctapproxy@sys-usb.service
